# installations/rc.local
#!/bin/sh -e

sudo jetson_clocks
sudo sh -c 'echo 1000 > /sys/module/usbcore/parameters/usbfs_memory_mb' &
exit 0
